# ubuntu-parental-control.service
[Unit]
Description=Ubuntu Parental Control Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/ubuntu-parental-control

# Set up iptables BEFORE starting the service
ExecStartPre=/opt/ubuntu-parental-control/setup_iptables.sh

# Start the service
ExecStart=/opt/ubuntu-parental-control/start_service.sh

# Clean up iptables when stopping
ExecStopPost=/usr/sbin/iptables -t nat -F PARENTAL_REDIRECT
ExecStopPost=/usr/sbin/iptables -t nat -D OUTPUT -j PARENTAL_REDIRECT
ExecStopPost=/usr/sbin/iptables -t nat -X PARENTAL_REDIRECT

# Restart on failure
Restart=on-failure
RestartSec=10

# Resource limits (prevent OOM)
MemoryMax=512M
MemoryHigh=384M

[Install]
WantedBy=multi-user.target
